# deploy/trading-dashboard.service
[Unit]
Description=Trading Bot Dashboard (Streamlit)
After=network.target

[Service]
Type=simple
User=trader
WorkingDirectory=/home/trader/trading-bot
EnvironmentFile=/home/trader/trading-bot/.env
ExecStart=/home/trader/.local/bin/uv run streamlit run dashboard.py --server.port 8501 --server.address 127.0.0.1 --server.headless true
Restart=always
RestartSec=10
StandardOutput=append:/home/trader/trading-bot/logs/dashboard.log
StandardError=append:/home/trader/trading-bot/logs/dashboard.log

[Install]
WantedBy=multi-user.target
