# deploy/trading-dashboard.service
[Unit]
Description=Trading Bot Dashboard (Streamlit)
After=network.target

[Service]
Type=simple
User=trader
WorkingDirectory=/home/trader/trading-bot
EnvironmentFile=/home/trader/trading-bot/.env
ExecStart=/home/trader/.local/bin/uv run streamlit run dashboard.py
Restart=always
RestartSec=10
AmbientCapabilities=CAP_NET_BIND_SERVICE
StandardOutput=append:/home/trader/trading-bot/logs/dashboard.log
StandardError=append:/home/trader/trading-bot/logs/dashboard.log

[Install]
WantedBy=multi-user.target
